# packages/system/.config/systemd/user/low-battery.timer
[Unit]
Description=Check batterie toutes les 3 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=3m
Persistent=true

[Install]
WantedBy=default.target
