# packages/system/.config/systemd/user/low-battery.service
[Unit]
Description=Alerte batterie faible

[Service]
Type=oneshot
ExecStart=%E/hypr/scripts/low_battery.sh
